# config/supervisor.conf
; SixBTC Supervisor Configuration
; Install: sudo cp config/supervisor.conf /etc/supervisor/conf.d/sixbtc.conf
;          sudo supervisorctl reread && sudo supervisorctl update
;
; Usage:
;   supervisorctl start sixbtc:*        # Start all
;   supervisorctl stop sixbtc:*         # Stop all
;   supervisorctl restart sixbtc:*      # Restart all
;   supervisorctl status sixbtc:*       # Status all

; ============================================================================
; ALL SIXBTC PROGRAMS
; ============================================================================

[program:generator]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/generator.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/generator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/generator.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/bitwolf",USER="bitwolf",PYTHONPATH="/home/bitwolf/sixbtc",PATH="/home/bitwolf/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
priority=100

[program:validator]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/validator.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/validator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/validator.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=200

[program:backtester]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/backtester.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/backtester.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/backtester.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=300

[program:classifier]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/classifier.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/classifier.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/classifier.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=400

[program:executor]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/executor.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/executor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/executor.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=500

[program:subaccount]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/subaccount_manager.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/subaccount.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/subaccount.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=600

[program:monitor]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/monitor.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=2
startretries=999999
stopwaitsecs=5
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/monitor.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=700

[program:scheduler]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/scheduler.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/bitwolf/sixbtc/logs/scheduler.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=800

; ============================================================================
; SINGLE GROUP - All programs under "sixbtc"
; ============================================================================

[group:sixbtc]
programs=generator,validator,backtester,classifier,executor,subaccount,monitor,scheduler
priority=999
