# config/supervisor.conf
; SixBTC Supervisor Configuration
; Install: sudo cp config/supervisor.conf /etc/supervisor/conf.d/sixbtc.conf
;          sudo supervisorctl reread && sudo supervisorctl update
;
; Usage:
;   supervisorctl start sixbtc:*        # Start all
;   supervisorctl stop sixbtc:*         # Stop all
;   supervisorctl restart sixbtc:*      # Restart all
;   supervisorctl status sixbtc:*       # Status all

; ============================================================================
; ALL SIXBTC PROGRAMS
; ============================================================================

[program:generator]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/generator.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/generator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=HOME="/home/bitwolf",USER="bitwolf",PYTHONPATH="/home/bitwolf/sixbtc",PATH="/home/bitwolf/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
priority=100

[program:validator]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/validator.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/validator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=200

[program:backtester]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/backtester.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/backtester.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=300

; DEPRECATED: classifier functionality moved to:
; - SCORER module (scoring) - integrated into BACKTESTER
; - ROTATOR module (selection + deployment)

[program:rotator]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/rotator.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/rotator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=400

[program:executor]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/executor.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/executor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=500

; DEPRECATED: subaccount manager functionality moved to ROTATOR module

[program:monitor]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/monitor.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=2
startretries=999999
stopwaitsecs=5
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=700

[program:scheduler]
command=/home/bitwolf/sixbtc/.venv/bin/python -u /home/bitwolf/sixbtc/src/processes/scheduler.py
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=unexpected
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=800

[program:metrics]
command=/home/bitwolf/sixbtc/.venv/bin/python -u -m src.metrics.collector
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/metrics.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=850

; ============================================================================
; WEB DASHBOARD
; ============================================================================

[program:api]
command=/home/bitwolf/sixbtc/.venv/bin/uvicorn src.api.main:app --host 0.0.0.0 --port 8080
directory=/home/bitwolf/sixbtc
user=bitwolf
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/home/bitwolf/sixbtc/logs/api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
environment=PYTHONPATH="/home/bitwolf/sixbtc"
priority=900

[program:frontend]
command=/usr/bin/npm run dev -- --host 0.0.0.0 --port 5173
directory=/home/bitwolf/sixbtc/web
user=bitwolf
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/home/bitwolf/sixbtc/logs/frontend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
priority=901

; ============================================================================
; SINGLE GROUP - All programs under "sixbtc"
; ============================================================================

[group:sixbtc]
programs=generator,validator,backtester,rotator,executor,monitor,scheduler,metrics,api,frontend
priority=999
